# systemd/pagefeed.socket
[Unit]
Description=Socket for pagefeed

[Socket]
SocketUser=www-data
ListenStream=/run/pagefeed.socket

[Install]
WantedBy=pagefeed.target
